# main/supervisord.conf
[program:node]
command = node server/app.js
directory = /home/dotcloud/current